# backend/systemd/stock-refresh.service
[Unit]
Description=Lynch Stock Screener - Nightly Data Refresh
After=network.target

[Service]
Type=oneshot
User=%u
WorkingDirectory=/home/user/lynch-stock-screener/backend
ExecStart=/usr/bin/python3 /home/user/lynch-stock-screener/backend/batch_refresh.py --config /home/user/lynch-stock-screener/backend/batch_config.json

# Environment
Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stock-refresh

# Resource limits (prevent runaway process)
MemoryLimit=2G
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
